# resources/namedmanager-www.cron
# NAMEDMANAGER CRONJOBS
#

# Synchronise any changes to Route53
*/1 * * * * www-data cd /opt/namedmanager/htdocs/include/cron/; php -q route53.php >> /var/log/namedmanager/cron-route53 2>&1

